# systemd/whatsapp-flask.service
[Unit]
Description=WhatsApp Scheduler Flask App
After=network.target whatsapp-driver.service
Requires=whatsapp-driver.service

[Service]
Type=simple
User=iedwards
WorkingDirectory=/home/iedwards/development/whatsapp_scheduler_openwa
ExecStart=/usr/bin/python3 app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
